# bt_speaker_addon/rootfs/etc/services.d/bt-speaker/finish
#!/usr/bin/with-contenv bashio
# Appelé quand le service s'arrête
bashio::log.warning "BT Speaker s'est arrêté."