# bt_speaker_addon/rootfs/etc/services.d/bt-speaker/finish
#!/usr/bin/with-contenv bashio
bashio::log.warning "BT Speaker s'est arrêté."