# bt_speaker_addon/rootfs/etc/services.d/bt-speaker/finish
#!/command/with-contenv bashio
bashio::log.warning "BT Speaker s'est arrêté."