# bt_speaker_addon/rootfs/etc/services.d/bt-speaker/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Démarrage BT Speaker..."

SPEAKER_MAC=$(bashio::config 'speaker_mac')
POLL_INTERVAL=$(bashio::config 'poll_interval')

mkdir -p /var/run/pulse
pulseaudio --system \
    --disallow-exit \
    --disallow-module-loading=false \
    --log-target=stderr \
    --exit-idle-time=-1 &

sleep 2
pactl load-module module-bluetooth-discover || true
pactl load-module module-bluetooth-policy   || true

exec /opt/venv/bin/python3 /app/main.py